# systemd/wap-kitty.service
[Unit]
Description=Generate kitty theme from wallpaper (M3 TonalSpot)

[Service]
Type=oneshot
ExecStart=%h/Downloads/wap-kitty/wap
